# etc/systemd/system/kiosk-config.service
[Unit]
Description=Apply per-panel kiosk config from /data
Before=kiosk.service
After=data.mount
Wants=data.mount

[Service]
Type=oneshot
ExecStart=/bin/sh -c '[ -f /data/poly-kiosk/config ] && cp -f /data/poly-kiosk/config /etc/default/tc8-kiosk || true'

[Install]
WantedBy=multi-user.target
