# etc/systemd/system/tc8-wired-mac.service
[Unit]
Description=Provision stable Polycom-OUI wired MAC (end0/lan) from SoC UID
Documentation=man:tc8-hwaddr
# Run before systemd-networkd configures the DSA ports so the interfaces come
# up (and DHCP) with the stable Polycom MAC, not the bootloader's random one.
Wants=network-pre.target
Before=network-pre.target
ConditionPathExists=/sys/devices/soc0/serial_number

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/tc8-hwaddr.sh net

[Install]
WantedBy=multi-user.target
