# etc/systemd/system/tc8-wired-mac.service
[Unit]
Description=Provision factory wired MAC (end0/lan) from the stock U-Boot env
Documentation=man:tc8-hwaddr
# Run before systemd-networkd configures the DSA ports so the interfaces come
# up (and DHCP) with the factory MAC, not the bootloader's random one.
Wants=network-pre.target
Before=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/tc8-hwaddr.sh net

[Install]
WantedBy=multi-user.target
